# docker/rootfs/etc/s6-overlay/s6-rc.d/mcp/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e

. /usr/bin/common.sh

cd /app/mcp || exit 1

# Talk to the local NPM API unless overridden.
export NPM_BASE_URL="${NPM_BASE_URL:-http://127.0.0.1:81/api}"
export MCP_PORT="${MCP_PORT:-3001}"
# Bind all interfaces so the endpoint is reachable outside the container.
export MCP_HOST="${MCP_HOST:-0.0.0.0}"

# The endpoint controls the whole proxy, so it must never come up open on a
# non-loopback bind. Without a bearer token the service parks itself.
if [ -z "${MCP_AUTH_TOKEN:-}" ]; then
	log_info 'MCP server disabled: set MCP_AUTH_TOKEN to enable the /mcp endpoint'
	exec sleep infinity
fi

if [ -z "${NPM_API_KEY:-}" ] && { [ -z "${NPM_IDENTITY:-}" ] || [ -z "${NPM_SECRET:-}" ]; }; then
	log_info 'MCP server disabled: set NPM_API_KEY or NPM_IDENTITY/NPM_SECRET'
	exec sleep infinity
fi

log_info 'Starting MCP server ...'

while :
do
	s6-setuidgid "$PUID:$PGID" bash -c "export HOME=$NPMHOME;node dist/index.js" || true
	sleep 2
done
